# FirmwarePackager/templates/scripts/init/sysv/S95-upgrade-recover.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          upgrade-recover
# Required-Start:    $local_fs
# Required-Stop:
# Default-Start:     S
# Default-Stop:
# Short-Description: Run recovery after upgrade
### END INIT INFO
case "$1" in
    start)
        /opt/upgrade/last/scripts/recover_boot.sh
        ;;
    *)
        echo "Usage: $0 start" >&2
        exit 1
        ;;
esac

exit 0
